# supervisord/fetch.ini
[program:gnpd_fetch]
command=sh /opt/app/gnpd/fetch/publish.sh
user=mapr
log_stdout=true
log_stderr=true
logfile=/opt/app/gnpd/supervisord/logs/fetch.log
logfile_maxbytes=1MB
logfile_backups=10
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=/opt/app/gnpd/supervisord/logs/fetch.stdout.log
